# frpc/rootfs/etc/services.d/frpc/run
#!/usr/bin/with-contenv bashio
set -e

CONFIG="/etc/frp/frpc.ini"
TEMPCONFIG="/tmp/frpc.ini"

# 根据 config.yaml 生成 frpc.ini
bashio::var.json \
    server_addr "$(bashio::config 'server_addr')" \
    server_port "$(bashio::config 'server_port')" \
    token "$(bashio::config 'token')" \
    tunnels "$(bashio::config 'tunnels')" \
    | tempio \
        -template /etc/frp/frpc.ini.gtpl \
        -out "${TEMPCONFIG}"

# 启动 frpc
exec /usr/bin/frpc -c "${TEMPCONFIG}"
